# files/etc/init/knot.conf
# knot - Knot DNS Server

description "Knot DNS Server"
author "Ondrej Sury <ondrej@debian.org>"

start on runlevel [2345]
stop on runlevel [016]

env CONFFILE=/etc/knot/knot.conf

respawn
pre-start script
  while ! /usr/local/bin/network_status.sh
  do
    :
  done
  exec /usr/lib/knot/prepare-environment $CONFFILE
end script
exec /usr/sbin/knotd -c $CONFFILE
